# examples/deploy/doggie-boot.service
[Unit]
Description=Custom Doggie Fast Boot
Wants=network-online.target
After=network-online.target sound.target

[Service]
Type=oneshot
RemainAfterExit=yes
User=root
WorkingDirectory=/home/matt/pidog
ExecStart=/bin/bash /home/matt/pidog/scripts/doggie_boot.sh
Environment=PYTHONUNBUFFERED=1
Environment=LOGNAME=matt
Environment=DOGGIE_REPO_DIR=/home/matt/pidog
Environment=DOGGIE_REPO_OWNER=matt
Environment=DOGGIE_BRANCH=main
Environment=DOGGIE_NETWORK_TIMEOUT=20
Environment=DOGGIE_PULL_TIMEOUT=25
Environment=DOGGIE_START_DELAY=8
Environment=DOGGIE_BOOT_RETRIES=3
Environment=DOGGIE_BOOT_RETRY_DELAY=5
TimeoutStartSec=70
TimeoutStopSec=15
Nice=-10
CPUWeight=800
IOWeight=800
IOSchedulingClass=best-effort
IOSchedulingPriority=0

[Install]
WantedBy=multi-user.target
